# discord-bot/sandbox-pal-dispatch-bot.service
[Unit]
Description=Sandbox Pal Dispatch Discord Bot
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
WorkingDirectory=WORKING_DIR
EnvironmentFile=CONFIG_PATH
ExecStart=WORKING_DIR/.venv/bin/python bot.py
Restart=on-failure
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=default.target
